# bin/systemd/iatistudio.service
[Unit]
Description=iati-studio

[Service]
WorkingDirectory=/home/oipa/IATI-STUDIO/app
User=oipa
Group=oipa
ExecStart=PORT=3000 /usr/bin/node index.js
Restart=always

[Install]
WantedBy=multi-user.target
